# linux/usb-mass-storage.service
[Unit]
Description=USB Mass Storage Gadget
After=local-fs.target
Before=getty@tty1.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/browser/usb_setup.sh
ExecStop=/home/browser/usb_disable.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target